# deploy/systemd/flashlight-api.service
[Unit]
Description=Flashlight Ratings API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=flashlight
Group=flashlight
WorkingDirectory=/opt/flashlight-ratings-go
EnvironmentFile=/etc/flashlight-ratings-go/api.env
ExecStart=/opt/flashlight-ratings-go/bin/api
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/opt/flashlight-ratings-go
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
