# deploy/systemd/flashlight-worker.service
[Unit]
Description=Flashlight Ratings Scheduler Worker
After=network-online.target flashlight-api.service
Wants=network-online.target

[Service]
Type=simple
User=flashlight
Group=flashlight
WorkingDirectory=/opt/flashlight-ratings-go
EnvironmentFile=/etc/flashlight-ratings-go/worker.env
ExecStart=/opt/flashlight-ratings-go/bin/worker
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/opt/flashlight-ratings-go
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
